# code/api-file-master/api-file-master/configuration/supervisord.d/dev/update_user_info_to_mongo.conf
[program:update_user_info_to_mongo]
process_name=%(program_name)s_%(process_num)02d
command=/usr/bin/php /data/www/api-version/artisan queue:work redis --queue=update_user_info_to_mongo --daemon --tries=10 --timeout=60
autostart=true
autorestart=true
numprocs=2
redirect_stderr=true
user=server
stdout_logfile=/var/log/supervisor/update_user_info_to_mongo.log
